# examples/launchd-plist.example.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  dashi-plugin-claude-code · launchd agent example для macOS

  Скопируйте в ~/Library/LaunchAgents/com.dashi-plugin.channel-<agent>.plist
  Замените:
    <you>      → ваш macOS username (echo $USER)
    myagent    → имя вашего агента (lowercase, без пробелов)

  Загрузка:
    launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/com.dashi-plugin.channel-myagent.plist

  См. docs/03-installation-macos.md для полной инструкции.
-->
<plist version="1.0">
<dict>
    <!-- Уникальный идентификатор agent'а в launchd namespace -->
    <key>Label</key>
    <string>com.dashi-plugin.channel-myagent</string>

    <!--
      ProgramArguments — массив argv процесса.
      Используем sh -c потому что:
      1. Нужно source secrets из channel.env (не хотим класть токены в plist —
         plist читается любым процессом этого user, env stays in memory)
      2. Нужно отправить Enter в tmux для прохождения welcome-промтов

      WARNING: путь к tmux зависит от вашей установки:
        Apple Silicon (M1/M2/M3/M4): /opt/homebrew/bin/tmux
        Intel Mac:                    /usr/local/bin/tmux
      Проверьте `which tmux` и подставьте.
    -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/sh</string>
        <string>-c</string>
        <string>set -a; . /Users/&lt;you&gt;/.claude-lab/myagent/secrets/channel.env; set +a; exec /opt/homebrew/bin/tmux new-session -d -s channel-myagent claude --dangerously-load-development-channels server:dashi-channel; sleep 6; /opt/homebrew/bin/tmux send-keys -t channel-myagent Enter; sleep 2; /opt/homebrew/bin/tmux send-keys -t channel-myagent Enter; wait</string>
    </array>

    <!--
      WorkingDirectory — КРИТИЧНО.
      Должен указывать ВНУТРЬ plugin/ каталога ВНУТРИ workspace агента.
      Claude Code от этого CWD ищет CLAUDE.md upward.
      См. docs/02-where-to-place-plugin.md
    -->
    <key>WorkingDirectory</key>
    <string>/Users/&lt;you&gt;/.claude-lab/myagent/.claude/dashi-plugin-claude-code/plugin</string>

    <!--
      EnvironmentVariables — переменные окружения для процесса.
      Только non-secret здесь. Bot token и т.п. — через source в ProgramArguments.
    -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/&lt;you&gt;</string>
        <!--
          PATH: bun, tmux, claude. Проверьте что claude binary доступен.
          Если Claude Code не в PATH — добавьте /Applications/Claude.app/Contents/...
        -->
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/Users/&lt;you&gt;/.bun/bin</string>
        <key>LANG</key>
        <string>en_US.UTF-8</string>
    </dict>

    <!-- Запуск при загрузке plist (вход в систему GUI user) -->
    <key>RunAtLoad</key>
    <true/>

    <!--
      KeepAlive: запускать снова если:
      - процесс упал с non-zero exit (Crashed)
      Если установить KeepAlive=true (булево, не dict) — будет рестарт всегда,
      что приведёт к loop когда welcome-промты не пройдены.
      Поэтому используем dict с conditional restart.
    -->
    <key>KeepAlive</key>
    <dict>
        <key>Crashed</key>
        <true/>
    </dict>

    <!-- Пауза 15 секунд между попытками рестарта (антиloop) -->
    <key>ThrottleInterval</key>
    <integer>15</integer>

    <!-- Логи. Создайте директорию заранее: mkdir -p ~/Library/Logs/dashi-plugin -->
    <key>StandardOutPath</key>
    <string>/Users/&lt;you&gt;/Library/Logs/dashi-plugin/channel-myagent.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/&lt;you&gt;/Library/Logs/dashi-plugin/channel-myagent.err.log</string>

    <!--
      Опционально: ProcessType: Interactive
      Гарантирует, что launchd не будет агрессивно ограничивать ресурсы
      (имеет смысл для interactive-приложений типа Claude Code, который читает TTY)
    -->
    <key>ProcessType</key>
    <string>Interactive</string>

    <!--
      Опционально: Nice — приоритет. 0 = default, -20 высокий, 20 низкий.
      Для агентов фоновых задач — оставьте default (закомментировано).
    -->
    <!--
    <key>Nice</key>
    <integer>1</integer>
    -->
</dict>
</plist>
