# examples/webhook-listener.service.example
[Unit]
Description=dashi-plugin webhook listener (gbrain swarm → claude -p)
After=network-online.target
Wants=network-online.target

[Service]
# Use a non-root service user. The user must own AGENT_WORKSPACE.
User=<service-user>
Group=<service-user>

EnvironmentFile=/etc/dashi-plugin/webhook.env

# HOME is needed by Claude Code for ~/.claude/CLAUDE.md and settings.json
Environment=HOME=/home/<service-user>
Environment=PATH=/usr/local/bin:/usr/bin:/bin

WorkingDirectory=/opt/dashi-plugin-webhook
ExecStart=/opt/dashi-plugin-webhook/.venv/bin/python /opt/dashi-plugin-webhook/listener.py

Restart=on-failure
RestartSec=5

# Hardening
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=read-only
ReadWritePaths=/var/log/dashi-plugin-webhook /home/<service-user>/.claude-lab

[Install]
WantedBy=multi-user.target
